# config_tools/mesh-network.service
[Unit]
Description=Mobile Ad-Hoc Mesh Network Service
Wants=NetworkManager-wait-online.service
After=NetworkManager-wait-online.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/mesh-network/mesh-config.conf

ExecStart=/usr/sbin/mesh-network.sh
ExecStop=/usr/sbin/mesh-network-stop.sh

[Install]
WantedBy=multi-user.target
