# config_tools/mesh-network.service
[Unit]
Description=Mesh Network Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/mesh-network.sh service
ExecStop=/usr/sbin/mesh-network-stop.sh
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3
RemainAfterExit=yes
EnvironmentFile=/etc/mesh-network/mesh-config.conf

[Install]
WantedBy=multi-user.target
